# docker/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:confd]
command=/usr/local/bin/confd -backend env
exitcodes=0


[program:nginx]
command=nginx -g "daemon off;"